# templates/task.conf.j2
[program:{{ name }}]
command={{ command }}
directory={{ directory }}
redirect_stderr={{ redirect_stderr }}
autostart={{ autostart }}
autorestart={{ autorestart }}
user={{ user }}
